# media-plugins/hyperion/files/hyperion.initd
#!/sbin/openrc-run

description="Hyperion - an open source ambient light software"
user="${RC_SVCNAME}"
group="${RC_SVCNAME}"
pidfile="/run/${RC_SVCNAME}.pid"

command="/usr/bin/hyperiond"
command_background="true"
command_args="--userdata /var/lib/${RC_SVCNAME}"

start_stop_daemon_args="--user ${user} --group ${group} --stdout-logger logger --stderr-logger logger"
